# deploy/systemd/vpn-sales-backup.service
[Unit]
Description=Create backup snapshot for VPN sales bot
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/vpn-sales-backup.sh
